# service/stable-diffusion-webui.service
[Unit]
Description=Stable Diffusion Web UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/stable-diffusion-webui
ExecStart=/usr/bin/python3 launch.py --listen --port 7860
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target