# files/stable-diffusion-webui.service
[Unit]
Description=Stable Diffusion Web UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/stable-diffusion-webui
ExecStart=/home/ubuntu/stable-diffusion-webui/webui.sh --listen --port 7860
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target